# root/etc/services.d/jdownloader/finish
#!/usr/bin/execlineb -S0

rm -f /opt/JDownloader/JDownloader.jar.*
rm -f /opt/JDownloader/JDownloader.pid